# veilhunter_run_key_suspicious_value.yml
title: VeilHunter - Run Key Value Points to Suspicious Path
id: be80e214-e344-453c-9012-cc758a25f0b7
status: experimental
description: >
  Detects a Run key value that references an executable in a suspicious
  staging location such as Temp, AppData, or ProgramData.
  Legitimate software rarely installs to these paths.
references:
  - https://attack.mitre.org/techniques/T1547/001/
author: VeilHunter
date: 2026-03-25
tags:
  - attack.persistence
  - attack.t1547.001
logsource:
  product: windows
  category: registry_set
detection:
  selection:
    TargetObject|contains:
      - '\CurrentVersion\Run\'
      - '\CurrentVersion\RunOnce\'
    Details|contains:
      - '\Temp\'
      - '\AppData\'
      - '\ProgramData\'
      - '\Users\Public\'
      - '\Downloads\'
  condition: selection
falsepositives:
  - Some portable or user-installed software
level: high
